# oem/eos/desktop-overlay/etc/skel/.config/autostart/plank.desktop
[Desktop Entry]
Name=Plank
Type=Application
Exec=/usr/bin/plank
Icon=plank
Comment=The main application dock
X-GNOME-Autostart-enabled=true