# oem/eos/live-overlay/etc/skel/.config/autostart/plank.desktop
[Desktop Entry]
Name=Plank
Type=Application
Exec=/usr/bin/plank
Icon=plank
Comment=A dock from Pantheon
X-GNOME-Autostart-enabled=true